# next-app/com.dungeonhookah.bot.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.dungeonhookah.bot</string>
    <key>ProgramArguments</key>
    <array>
        <string>/Users/nikolajmisin/Documents/WebApp/next-app/start-bot-forever.sh</string>
    </array>
    <key>WorkingDirectory</key>
    <string>/Users/nikolajmisin/Documents/WebApp/next-app</string>
    <key>RunAtLoad</key>
    <true/>
    <key>KeepAlive</key>
    <true/>
    <key>StandardOutPath</key>
    <string>/Users/nikolajmisin/Documents/WebApp/next-app/bot-launchd.log</string>
    <key>StandardErrorPath</key>
    <string>/Users/nikolajmisin/Documents/WebApp/next-app/bot-launchd.log</string>
</dict>
</plist>
