# next-app/dungeonhookah-bot.service
[Unit]
Description=Dungeon Hookah Telegram Bot
After=network.target

[Service]
Type=simple
User=nikolajmisin
WorkingDirectory=/Users/nikolajmisin/Documents/WebApp/next-app
ExecStart=/Users/nikolajmisin/Documents/WebApp/next-app/start-bot-forever.sh
Restart=always
RestartSec=10
StandardOutput=append:/Users/nikolajmisin/Documents/WebApp/next-app/bot-service.log
StandardError=append:/Users/nikolajmisin/Documents/WebApp/next-app/bot-service.log

[Install]
WantedBy=multi-user.target
